# rules/local_rules.xml
<!-- Larkspur Retail Group — Custom Wazuh Detection Rules -->
<!-- CA1 Assessment | B9CY110 | Samuel Ogunlusi | 20086108 -->
<!-- MITRE ATT&CK mappings: T1110, T1059.001, T1053.003, T1078 -->

<group name="larkspur,custom,">

  <!-- Rule 100001: Windows brute-force logon failures (T1110) -->
  <rule id="100001" level="10" frequency="5" timeframe="120">
    <if_matched_sid>60122</if_matched_sid>
    <mitre>
      <id>T1110</id>
    </mitre>
    <description>Multiple Windows logon failures detected — possible brute-force attack</description>
    <group>authentication_failures,pci_dss_10.2.4,pci_dss_10.2.5,gpg13_7.1,</group>
  </rule>

  <!-- Rule 100002: SSH brute-force (T1110) -->
  <rule id="100002" level="10" frequency="5" timeframe="120">
    <if_matched_sid>5710</if_matched_sid>
    <mitre>
      <id>T1110</id>
    </mitre>
    <description>Multiple SSH login failures — possible brute-force attack against Linux endpoint</description>
    <group>authentication_failures,pci_dss_10.2.4,pci_dss_10.2.5,</group>
  </rule>

  <!-- Rule 100003: Suspicious PowerShell execution (T1059.001) -->
  <rule id="100003" level="12">
    <if_sid>91801,91802,91803</if_sid>
    <mitre>
      <id>T1059.001</id>
    </mitre>
    <description>Suspicious PowerShell script execution detected via Script Block Logging</description>
    <group>attack,pci_dss_10.6.1,</group>
  </rule>

  <!-- Rule 100004: Cron persistence via FIM (T1053.003) -->
  <rule id="100004" level="8">
    <if_sid>2930,2931</if_sid>
    <field name="file">/etc/cron</field>
    <mitre>
      <id>T1053.003</id>
    </mitre>
    <description>Cron configuration file modified — possible persistence mechanism installed</description>
    <group>syscheck,pci_dss_11.5,</group>
  </rule>

  <!-- Rule 100005: Non-existent user SSH reconnaissance (T1078) -->
  <rule id="100005" level="8" frequency="3" timeframe="60">
    <if_matched_sid>5710</if_matched_sid>
    <match>Invalid user</match>
    <mitre>
      <id>T1078</id>
    </mitre>
    <description>Repeated SSH login attempts using non-existent usernames — reconnaissance activity</description>
    <group>authentication_failures,pci_dss_10.2.4,</group>
  </rule>

</group>
